# bin/harchiverd
#!/bin/sh
#http://blog.scphillips.com/2013/07/getting-a-python-script-to-run-in-the-background-as-a-service-on-boot/
 
### BEGIN INIT INFO
# Provides:		harchiverd
# Short-Description:	HAR-archiving daemon.
# Description:		Stores HAR records generated via a webservice in
#			WARC files.
### END INIT INFO
 
# Change the next 3 lines to suit where you install your script and what you want to call it
DIR="/usr/local/bin"
DAEMON=$DIR/harchiver.py
DAEMON_NAME=harchiverd
NUM_DAEMONS=4
 
# This next line determines what user the script runs as.
# Root generally not recommended but necessary if you are using the Raspberry Pi GPIO from Python.
DAEMON_USER=heritrix
 
# The process ID of the script when it runs is stored here:
PIDFILE=/var/run/$DAEMON_NAME/$DAEMON_NAME.pid
 
. /lib/lsb/init-functions
 
do_start () {
	if [[ $NUM_DAEMONS -gt 1 ]]
	then
		for i in $(eval echo "{1..$NUM_DAEMONS}")
		do
			su - $DAEMON_USER -c "$DAEMON start /var/run/$DAEMON_NAME/$DAEMON_NAME-$i.pid"
		done
	else
		su - $DAEMON_USER -c "$DAEMON start"
	fi
}
do_status () {
	if [[ $NUM_DAEMONS -gt 1 ]]
	then
		for i in $(eval echo "{1..$NUM_DAEMONS}")
		do
			echo -n "$i "
			if [[ -f "/var/run/$DAEMON_NAME/$DAEMON_NAME-$i.pid" ]]
			then
				if kill -0 $(cat "/var/run/$DAEMON_NAME/$DAEMON_NAME-$i.pid") 1> /dev/null 2> /dev/null
				then
					echo "$DAEMON_NAME is running."
				else
					echo "$DAEMON_NAME is not running."
				fi
			else
				echo "$DAEMON_NAME pidfile not found."
			fi
		done
	else
		if [[ -f $PIDFILE ]]
		then
			if kill -0 $(cat $PIDFILE) 1> /dev/null 2> /dev/null
			then
				echo "$DAEMON_NAME is running."
			else
				echo "$DAEMON_NAME is not running."
			fi
		else
			echo "$DAEMON_NAME pidfile not found."
		fi
	fi
}
do_stop () {
	if [[ $NUM_DAEMONS -gt 1 ]]
	then
		for i in $(eval echo "{1..$NUM_DAEMONS}")
		do
			su - $DAEMON_USER -c "$DAEMON stop /var/run/$DAEMON_NAME/$DAEMON_NAME-$i.pid"
		done
	else
		su - $DAEMON_USER -c "$DAEMON stop"
	fi
}
 
case "$1" in
 
	start|stop)
		do_${1}
		;;

	status)
		do_status
		;;
 
	restart)
		do_stop
		do_start
		;;
 
	*)
		echo "Usage: /etc/init.d/$DAEMON_NAME {start|stop|status|restart}"
		exit 1
		;;
 
esac
exit 0
